# etc/profile.d/97.rotifer.sh
append_directory_to_path /home/kaihami/bin
add_to_environment_variable PYTHONPATH first /home/kaihami/mymodules/rotifer/lib
