# etc/profile.d/98.rotifer.bash_completion.sh
eval "$(register-python-argcomplete conseq)"

eval "$(register-python-argcomplete pfetch)"

eval "$(register-python-argcomplete trimal2msa)"

eval "$(register-python-argcomplete acc2operon)"
eval "$(register-python-argcomplete /home/kaihami/bin/acc2operon)"
eval "$(register-python-argcomplete /home/kaihami/mymodules/rotifer/bin/acc2operon.py)"

# eval "$(register-python-argcomplete ~kaihami/bin/acc2operon -s bash)" not working
# eval "$(register-python-argcomplete ~kaihami/mymodules/rotifer/bin/acc2operon.py -s bash)"

eval "$(register-python-argcomplete rconfig)"

eval "$(register-python-argcomplete rmsd_comparer)"

eval "$(register-python-argcomplete uniprot2gi)"

eval "$(register-python-argcomplete rnexplorer)"
